# bitcoind/conf/supervisord.conf
[supervisord]
nodaemon         = true
logfile          = /logs/supervisord.log
childlogdir      = /logs
logfile_maxbytes = 50MB
logfile_backups  = 5
loglevel         = warn
pidfile          = /var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[rpcinterface:ccunrpc]
supervisor.rpcinterface_factory = ccunrpc.bitcoind:make_main_rpcinterface
daemon_config = /config/bitcoin.conf

[inet_http_server]
port = 0.0.0.0:9001
#username = dc4ecrfref34fewcf
#password = fevgrebrtb5revfer34rfv

[include]
files = supervisord_custom.conf

[program:bitcoind]
command = bitcoind -conf=/config/bitcoin.conf
directory = /data
environment=HOME="/home/bitcoin",USER="bitcoin"
user = bitcoin
stopsignal=TERM
stopwaitsecs=300
autostart=true
autorestart=false
redirect_stderr=True
logfile_maxbytes = 50MB
loglevel         = warn
stdout_logfile = /logs/bitcoind_supervisor.log
